# code/puppet/modules/pub_sub_proxy/files/pub_sub_frontend.conf
# frontend - an example ZeroMQ express.js app

description "ZeroMQ Publisher express.js app"

start on pubsub_start
start on started pub_sub_proxy

stop on pubsub_stop
stop on stopped pub_sub_proxy

# respawn the job up to 10 times within a 5 second period.
# If the job exceeds these values, it will be stopped and
# marked as failed.
# respawn
# respawn limit 10 5

# Send Standard Out and Error to a log file in /var/log/upstart
console log

env NODE_ENV=production
env NODE_PATH=/usr/lib/nodejs:/usr/lib/node_modules:/usr/share/javascript

script 
  export APP_HOME=/vagrant/pub_sub_proxy/frontend
  exec /usr/bin/node $APP_HOME/app.js
end script

