# code/puppet/modules/pub_sub_proxy/files/sub_pub_worker.conf
# pub_sub_worker - a pub_sub worker for zeromq

description "ZeroMQ worker for pubsub"

start on pubsub_start
stop on pubsub_stop

# respawn the job up to 10 times within a 5 second period.
# If the job exceeds these values, it will be stopped and
# marked as failed.
respawn
respawn limit 10 5

# Send Standard Out and Error to a log file in /var/log/upstart
console log

exec /usr/bin/node /vagrant/pub_sub_proxy/worker/app.js
